# scripts/profile.d/801.fzf.sh
export FZF_HOME="$HOME"
[ -f $FZF_HOME/.fzf.zsh ] && source $FZF_HOME/.fzf.zsh
