# systemd/user/gnome-screensaver.service
[Unit]
Description=Gnome Screensaver
After=dbus.service

[Service]
ExecStart=/usr/bin/gnome-screensaver

[Install]
WantedBy=default.target
